# contrib/init/org.bitcoin.bitcoind.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
<<<<<<< HEAD
	<string>org.bitcoin.bitcoind</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/bitcoind</string>
=======
	<string>org.bitnova.bitnovad</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/bitnovad</string>
>>>>>>> 5360f2baff (Initialized BitNova project)
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
